# example/conf/supervisor.d/application.ini
[program:simple_script]
command=sh /root/eventhandler/simple_script.sh
autostart=true
autorestart=true
